# price-tracker.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Redirect job output to PID 1 so logs are visible through standard Docker logging.
0 */12 * * * root cd /app && /usr/local/bin/python -m scraper.main -c scraper/config.yml >> /proc/1/fd/1 2>> /proc/1/fd/2
